# etc/launchd/dev.yamen.vsc-tunnel.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.yamen.vsc-tunnel</string>
    <key>ProgramArguments</key>
    <array>
        <string>@AI_CONFIG_ROOT@/bin/vsc-tunnel</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>@HOME@/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>15</integer>
    <key>StandardOutPath</key>
    <string>@HOME@/Library/Logs/vsc-tunnel.log</string>
    <key>StandardErrorPath</key>
    <string>@HOME@/Library/Logs/vsc-tunnel.log</string>
</dict>
</plist>
